# fc.conf
[program:fc]
directory=/path/to/fc
command=/path/to/gunicorn fc.wsgi --bind=0.0.0.0:8080
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
environment=
